# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/signal/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── Log Follower ─────────────────────────────────────────────────────────────
[program:log_follower]
command=python -m ingestion.log_follower
directory=/opt/signal
user=signal
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=10
stdout_logfile=/var/log/signal/log_follower.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/var/log/signal/log_follower.err
stderr_logfile_maxbytes=10MB

; ── Enrichment Worker ────────────────────────────────────────────────────────
[program:enrichment_worker]
command=python -m enrichment.worker
directory=/opt/signal
user=signal
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=10
stdout_logfile=/var/log/signal/enrichment.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/var/log/signal/enrichment.err
stderr_logfile_maxbytes=10MB

; ── API Server ────────────────────────────────────────────────────────────────
[program:api]
command=uvicorn api.main:app --host 127.0.0.1 --port 8000 --workers 2
directory=/opt/signal
user=signal
autostart=true
autorestart=true
startsecs=3
startretries=10
stopwaitsecs=10
stdout_logfile=/var/log/signal/api.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile=/var/log/signal/api.err
stderr_logfile_maxbytes=10MB

; ── Monitor ──────────────────────────────────────────────────────────────────
[program:monitor]
command=python scripts/monitor.py
directory=/opt/signal
user=signal
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/dev/null
stderr_logfile=/var/log/signal/monitor.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5

; ── Process group ─────────────────────────────────────────────────────────────
[group:signal]
programs=log_follower,enrichment_worker,api,monitor
